# archlinux/backlight.service
[Unit]
Description = change permission of /sys/class/backlight/XXX/brightness

[Service]
Type=simple
ExecStart=/usr/bin/sleep 3
ExecStartPost=/usr/bin/chmod 666 /sys/class/backlight/nvidia_0/brightness
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target

# vim: ft=systemd
