# archlinux/rc.local
#!/usr/bin/env bash

echo 'power on' | /usr/bin/bluetoothctl

nohup socat TCP4-LISTEN:10086,reuseaddr,fork TCP6:[2001:19f0:9002:3ad:5400:01ff:fe3e:99a7]:8388 > /dev/null &

nohup socat UDP4-LISTEN:10086,reuseaddr,fork UDP6:[2001:19f0:9002:3ad:5400:01ff:fe3e:99a7]:8388 > /dev/null &

